# deploy/systemd/tradeagent-frontend.service
[Unit]
Description=tradeAgent Next.js frontend
After=network-online.target tradeagent-backend.service
Wants=network-online.target

[Service]
Type=simple
User=tradeagent
Group=tradeagent
WorkingDirectory=/opt/tradeagent/frontend
Environment="NODE_ENV=production"
Environment="BACKEND_URL=http://127.0.0.1:8000"
Environment="PORT=3000"
ExecStart=/opt/tradeagent/frontend/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
